# atomic/templates/mdg-backup.timer
[Unit]
Description=Run mdg backup everyday at 21:21

[Timer]
OnCalendar=*-*-* 21:21:00